# src/proxmox_srvbackup/adapters/config/systemd/proxmox-srvbackup-backup.service
[Unit]
Description=Proxmox Server Backup (Pull)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/proxmox-srvbackup backup
StandardOutput=journal
StandardError=journal
